# conf/etc/s6/services/01-mariadb/run
#!/bin/sh
NAME="MARIADB"
COMMAND="su - mysql -s /bin/sh -c 'mysqld'"
WAIT="/tmp/mariadb"

exec service-handler "$NAME" "$COMMAND" "$WAIT"
