# conf/etc/s6/services/01-mariadb/run
#!/bin/sh

service-handler "MARIADB" "su - mysql -s /bin/sh -c 'mysqld'" "/tmp/mariadb"